# configs/conf.d/main.conf
; ================================
;  celery bot supervisor
; ================================

[program:celerybot]
command=python3 /home/denis/PycharmProjects/Bots/LotoBot/bot.py --loglevel=INFO
directory=/home/denis/PycharmProjects/Bots/LotoBot

user=denis
numprocs=1
stdout_logfile=/var/log/celery/bot.log
stderr_logfile=/var/log/celery/bot.log
autostart=true
autorestart=true
; startsecs=10

; Causes supervisor to send the termination signal (SIGTERM) to the whole process group.
stopasgroup=true

; if rabbitmq is supervised, set its priority higher
; so it starts first
priority=999
